# initramfs/etc/inittab
::restart:/sbin/init
::shutdown:echo "Sync all file buffers..."
::shutdown:sync
::shutdown:echo "Unmount all filesystems..."
::shutdown:umount -a -r
::ctrlaltdel:/sbin/reboot
::once:cat /etc/motd
::respawn:/bin/cttyhack /bin/sh
tty2::once:cat /etc/motd
tty2::respawn:/bin/sh
tty3::once:cat /etc/motd
tty3::respawn:/bin/sh
tty4::once:cat /etc/motd
tty4::respawn:/bin/sh
